# resources/docs/platform/nginx-nodejs/config.sample/crontab
* * * * * cd /var/www && node schedule-script.js >> /dev/stdout 2>&1

